# serverConfig/db/crontab
*/30 * * * * /app/serverConfig/db/dbBackups.sh >/dev/null 2>&1
